# scripts/nanobot-gateway-user.service
[Unit]
Description=nanobot AI Assistant Gateway (User Service)
Documentation=https://github.com/HKUDS/nanobot
After=network.target
Wants=network.target

[Service]
Type=simple

# Working directory
WorkingDirectory=/home/nn/.nanobot/workspace/nano

# Environment
Environment="PATH=/home/nn/.nanobot/workspace/nano/.venv/bin:/home/nn/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="NANOBOT_WORKSPACE=/home/nn/.nanobot/workspace"
Environment="PYTHONUNBUFFERED=1"

# Main process with auto-restart
ExecStart=/home/nn/.local/bin/uv run nanobot gateway
Restart=always
RestartSec=5

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

# Security hardening (user-level compatible)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

# Allow write access to workspace and config
ReadWritePaths=/home/nn/.nanobot/workspace /home/nn/.nanobot/config

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nanobot-gateway

[Install]
WantedBy=default.target
